# docker/supervisord.conf
; =============================================================================
; supervisord.conf — MLB Prediction System process manager
;
; Manages two long-running processes inside the container:
;   mlb-predict-server  — uvicorn web server (FastAPI dashboard)
;   cron               — supercronic running docker/crontab
;
; supervisorctl can be used from inside the container to inspect or restart:
;   supervisorctl status
;   supervisorctl restart mlb-predict-server
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

; Unix socket for supervisorctl communication (used by the cron scripts)
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ---------------------------------------------------------------------------
; Web server
;
; Binds to 0.0.0.0 (all interfaces) inside the container so that Docker port
; forwarding works correctly.  On the host it is exposed on PORT (default 30087).
;
; %(ENV_MODEL)s and %(ENV_PORT)s are substituted from container env vars.
; ---------------------------------------------------------------------------
[program:mlb-predict-server]
command=python /app/scripts/serve.py
    --host 0.0.0.0
    --port %(ENV_PORT)s
    --model %(ENV_MODEL)s
directory=/app
environment=MLB_PREDICT_MODEL_TYPE=%(ENV_MODEL)s
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/app/logs/server.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/server.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5

; ---------------------------------------------------------------------------
; Cron daemon
;
; supercronic reads docker/crontab and executes scheduled tasks.  It logs
; each job's stdout/stderr to /app/logs/cron.log.
;
; Environment variables are explicitly forwarded so that cron scripts have
; access to MODEL and PORT regardless of supervisord's env inheritance.
; ---------------------------------------------------------------------------
[program:cron]
command=supercronic /app/docker/crontab
directory=/app
environment=MODEL=%(ENV_MODEL)s,PORT=%(ENV_PORT)s
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/cron.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/cron.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
